# conf/options/tpl_watcher.cron
*/5 * * * * (cd ${PRJ_ROOT}; sh bin/watcher.sh > /dev/null 2>&1 )
